# projects/computer-setup/2026-06-26/99-monitor-hotplug.rules
# 显示器热插拔自动检测 — 触发 monitor-hotplug.sh
# 当 DP/eDP/HDMI 接口状态改变时执行
ACTION=="change", SUBSYSTEM=="drm", RUN+="/usr/local/bin/monitor-hotplug.sh"
